# templates/etc/systemd/system/ansible-test-yamllint.timer.j2
# {{ ansible_managed }}
# ansibleguy.cicd_tester

[Unit]
Description=Timer to run {{ test }} of repository '{{ repo_name }}'

[Timer]
OnCalendar={{ tester.trigger.yamllint.day }} {{ trigger_time_yamllint }}:00
Persistent=false
WakeSystem=false

[Install]
WantedBy=multi-user.target


